# packaging/deb/postinst
#!/bin/sh
set -e

case "$1" in
  configure)
    # Defensive: ensure scripts are executable even if a tarball-based
    # build path stripped the mode bits.
    if [ -d /opt/Phantom-OS-KubernetesOptions/scripts ]; then
      find /opt/Phantom-OS-KubernetesOptions/scripts \
        -type f \( -name '*.sh' -o -name '*.py' \) \
        -exec chmod 0755 {} +
    fi

    # GPU/NVMAP host metric collectors (SOF-1167 #6) — install + enable straight
    # from the deb (deb-install-only) instead of scripts/install-gaia-host-services.sh,
    # so a plain `dpkg -i` brings them up. They write node-exporter textfiles to
    # /var/lib/gaia/node-textfile, which the k0s node-exporter DaemonSet serves.
    # (The docker-run gaia-tools RAG units stay in install-gaia-host-services.sh.)
    GAIA_HS=/opt/Phantom-OS-KubernetesOptions/host-services/gaia
    if command -v systemctl >/dev/null 2>&1 && [ -d "$GAIA_HS" ]; then
      mkdir -p /var/lib/gaia/node-textfile
      for s in gaia-gpu-metrics.sh jetson-nvmap-mem-textfile.sh; do
        [ -f "$GAIA_HS/$s" ] && install -m 0755 "$GAIA_HS/$s" "/usr/local/bin/$s"
      done
      for u in gaia-gpu-metrics.service gaia-nvmap-metrics.service; do
        [ -f "$GAIA_HS/$u" ] && install -m 0644 "$GAIA_HS/$u" "/etc/systemd/system/$u"
      done
      systemctl daemon-reload || true
      systemctl enable --now gaia-gpu-metrics.service gaia-nvmap-metrics.service >/dev/null 2>&1 || true
    fi

    cat <<'EOF'

phantomos-k0s installed at /opt/Phantom-OS-KubernetesOptions

Next steps on a fresh robot (run from the directory holding the image
bundle, i.e. the .deb + .tar.zst pair):
  sudo bash /opt/Phantom-OS-KubernetesOptions/scripts/install-image-bundle.sh ./
  sudo bash /opt/Phantom-OS-KubernetesOptions/scripts/configure-host.sh
  sudo bash /opt/Phantom-OS-KubernetesOptions/scripts/bootstrap-robot.sh

Per-host config lives at /etc/phantomos/host-config.yaml and is NOT
managed by this package -- it is generated by configure-host.sh.

EOF

    # After unpack, /opt/Phantom-OS-KubernetesOptions is a fresh git
    # checkout (RFC 0006). If it shows uncommitted modifications, the
    # operator hand-modified tracked files between installs. ArgoCD
    # reconciles against /opt/.../HEAD, ignoring those modifications.
    # Emit a soft warning so the operator knows their edits are inert
    # (without failing the install).
    INSTALL_TREE=/opt/Phantom-OS-KubernetesOptions
    if [ -d "$INSTALL_TREE/.git" ] && command -v git >/dev/null 2>&1; then
        if [ -n "$(git -C "$INSTALL_TREE" status --porcelain 2>/dev/null)" ]; then
            echo "warning: $INSTALL_TREE has local modifications after install" >&2
            echo "  ArgoCD reconciles against /opt/.../HEAD -- uncommitted changes are inert" >&2
            echo "  see git -C $INSTALL_TREE status for the modified files" >&2
        fi
    fi
    ;;
  abort-upgrade|abort-remove|abort-deconfigure)
    ;;
  *)
    echo "postinst called with unknown argument \`$1'" >&2
    exit 1
    ;;
esac

exit 0
